# build/live-build/config/includes.chroot/etc/systemd/system/greetd.service.d/10-seatd.conf
[Unit]
Wants=seatd.service
After=seatd.service
