# build/live-build/config/includes.chroot/etc/xdg/autostart/agenos-installer.desktop
[Desktop Entry]
Type=Application
Name=AgenOS Installer Autostart
Comment=Launch AgenOS installer automatically in the live session
Exec=/usr/local/bin/agenos-autostart-installer
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
