# systemd/kiosk.service
[Unit]
Description=AgenOS Kiosk Session
After=graphical.target ui.service
Wants=ui.service

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/cage -- /usr/local/bin/agenos-installer-app
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
